# docker-build/supervisor/conf.d/freeswitch.conf
[program:freeswitch]
command=/usr/local/freeswitch/bin/freeswitch -nonat
directory=/usr/local/freeswitch
user=root              
autostart=true                
autorestart=true
redirect_stderr=true
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10